# scripts/arctic-media.service
[Unit]
Description=Arctic Media Server
After=network.target
Wants=network.target

[Service]
Type=exec
User=arctic
Group=arctic
WorkingDirectory=/opt/arctic-media
Environment=PATH=/opt/arctic-media/.venv/bin
Environment=ARCTIC_MEDIA_ROOT=/opt/arctic-media/data
Environment=ARCTIC_TRANSCODE_DIR=/opt/arctic-media/transcode
Environment=FFMPEG_BIN=ffmpeg
Environment=FFMPEG_PRESET=veryfast
Environment=HLS_SEG_DUR=2.0
ExecStart=/opt/arctic-media/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=arctic-media

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/arctic-media/data /opt/arctic-media/transcode /opt/arctic-media/arctic.db

[Install]
WantedBy=multi-user.target
